# deploy/systemd/quantpilot-backup.service
[Unit]
Description=Create a verified QuantPilot release backup
After=network-online.target

[Service]
Type=oneshot
User=quantpilot
Group=quantpilot
WorkingDirectory=/opt/quantpilot/current
EnvironmentFile=/etc/quantpilot/quantpilot.env
Environment=QUANTPILOT_BACKUP_ROOT=/var/backups/quantpilot
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run db:backup:release
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/backups/quantpilot /var/lib/quantpilot
