# deploy/systemd/quantpilot-auth-cleanup.timer
[Unit]
Description=Run QuantPilot auth retention cleanup daily

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m
Unit=quantpilot-auth-cleanup.service

[Install]
WantedBy=timers.target
